# deploy/systemd/heidi-tender-certbot-renew.service
[Unit]
Description=Renew Heidi Tender Let's Encrypt certificates
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/default/heidi-tender-prod
ExecStart=/bin/sh -lc 'cd "$PROJECT_DIR" && /usr/bin/docker compose --env-file "$COMPOSE_ENV_FILE" -f "$COMPOSE_FILE" exec -T certbot certbot renew --webroot -w /var/www/certbot --quiet'
ExecStartPost=/bin/sh -lc 'cd "$PROJECT_DIR" && /usr/bin/docker compose --env-file "$COMPOSE_ENV_FILE" -f "$COMPOSE_FILE" exec -T nginx /usr/local/bin/render-nginx-config.sh && /usr/bin/docker compose --env-file "$COMPOSE_ENV_FILE" -f "$COMPOSE_FILE" exec -T nginx nginx -s reload'
